# supervisor.conf
[program:mesa]
command = /vagrant/run.sh
autostart=true
autorestart=true
user = vagrant
stdout_logfile = /var/log/mesa.log
redirect_stderr = true
